# ansible/roles/k3s_master/templates/etcd-defrag.timer.j2
[Unit]
Description=Weekly etcd defragmentation

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
